# roles/mygpo/templates/mygpo.socket
[Unit]
Description=mygpo socket

[Socket]
ListenStream=/run/mygpo.sock

[Install]
WantedBy=sockets.target